# tests/Console/data/test.cron
* * * * * test:task
0 */1 * * * test:command
3 * * * * test:bungled